# crookery/crookery.service.template
[Unit]
Description=Crookery Dashboard
After=network.target

[Service]
Type=simple
WorkingDirectory=CROOKERY_BACKEND_DIR
Environment=PYTHONPATH=.
Environment=CROOKERY_HOST=0.0.0.0
Environment=CROOKERY_PORT=4123
EnvironmentFile=-CROOKERY_BACKEND_DIR/.env
ExecStart=CROOKERY_UVICORN main:app --host $CROOKERY_HOST --port $CROOKERY_PORT
Restart=on-failure
RestartSec=5

[Install]
WantedBy=default.target
